# root/etc/s6-overlay/s6-rc.d/init-pmm-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cp /app/pmm/config/config.yml.template /config

# permissions
lsiown -R abc:abc \
    /config
